# bambuddy-beta/rootfs/etc/services.d/bambuddy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here
declare debug
debug=$(bashio::config 'debug')

bashio::log.info 'Link data dir'
mkdir -p /config/data
mkdir -p /config/logs
rm -rf /app/data
rm -rf /app/logs
ln -s /config/data /app/data
ln -s /config/logs /app/logs

cd /app

if [ "${debug}" = "true" ]; then
  bashio::log.info 'Enabling debug mode'
  export DEBUG=true
fi

export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

# Trusted Frame Origins
if bashio::config.exists 'trusted_frame_origins'; then
  ORIGINS=$(bashio::config 'trusted_frame_origins | join(",")')
  if [ -n "${ORIGINS}" ]; then
    bashio::log.info "Setting TRUSTED_FRAME_ORIGINS: ${ORIGINS}"
    export TRUSTED_FRAME_ORIGINS="${ORIGINS}"
  fi
fi

# External storage (Share / Media)
EXTERNAL_ROOTS=""
if bashio::config.true 'enable_share'; then
  bashio::log.info "Enabling /share as external root"
  EXTERNAL_ROOTS="/share"
fi
if bashio::config.true 'enable_media'; then
  bashio::log.info "Enabling /media as external root"
  if [ -n "${EXTERNAL_ROOTS}" ]; then
    EXTERNAL_ROOTS="${EXTERNAL_ROOTS}:/media"
  else
    EXTERNAL_ROOTS="/media"
  fi
fi
if [ -n "${EXTERNAL_ROOTS}" ]; then
  bashio::log.info "Setting BAMBUDDY_EXTERNAL_ROOTS: ${EXTERNAL_ROOTS}"
  export BAMBUDDY_EXTERNAL_ROOTS="${EXTERNAL_ROOTS}"
fi

# Self-signed CA certificate support
if bashio::config.true 'use_system_trust_store'; then
  CERTFILE=$(bashio::config 'certfile')
  CERT_SRC="/config/${CERTFILE}"
  if [ -n "${CERTFILE}" ] && [ -f "${CERT_SRC}" ]; then
    bashio::log.info "Installing CA certificate: ${CERTFILE}"
    cp "${CERT_SRC}" "/usr/local/share/ca-certificates/${CERTFILE}"
    update-ca-certificates
    export USE_SYSTEM_TRUST_STORE=true
    bashio::log.info "System trust store enabled"
  else
    bashio::log.warning "use_system_trust_store is enabled but certificate file not found: ${CERT_SRC}"
    bashio::log.warning "Place your CA certificate (.crt) in the addon_configs folder and set certfile to the correct filename"
  fi
fi

bashio::log.info 'Start Bambuddy ...'
uvicorn backend.app.main:app --host 0.0.0.0 --port 8000 --loop asyncio
